# toypad.service.dist
[Unit]
Description=ToyPad

[Service]
Type=simple
ExecStart=${CPD}/main.py

[Install]
WantedBy=multi-user.target
